# deployment/homegym.service
# HomeGym Gunicorn Systemd Service
# ==================================
#
# Installation:
#   1. Pfade und Umgebungsvariablen anpassen (siehe unten)
#   2. sudo cp deployment/homegym.service /etc/systemd/system/
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable homegym
#   5. sudo systemctl start homegym
#
# Status prüfen:
#   sudo systemctl status homegym
#
# Logs anschauen:
#   sudo journalctl -u homegym -f

[Unit]
Description=HomeGym Fitness Tracker Gunicorn Daemon
After=network.target mysql.service
Requires=mysql.service

[Service]
Type=simple
User=your_username          # ANPASSEN: Linux User
Group=your_group            # ANPASSEN: Linux Group (z.B. www-data, psaserv)

# ANPASSEN: Pfad zu deinem Projekt
WorkingDirectory=/path/to/your/project

# ANPASSEN: Pfad zu venv
Environment="PATH=/path/to/your/project/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=config.settings"

# WICHTIG: Secrets NICHT hier eintragen – ausschliesslich via .env!
# Die .env-Datei muss im WorkingDirectory liegen (siehe .env.example).
# Benötigte Variablen (GENAU diese Namen – settings.py ist case-sensitive):
#   SECRET_KEY        (NICHT DJANGO_SECRET_KEY – wird ignoriert!)
#   DB_ENGINE, DB_NAME, DB_USER, DB_PASSWORD, DB_HOST, DB_PORT
#   ALLOWED_HOSTS
#
# WARNUNG: Environment="DJANGO_SECRET_KEY=..." funktioniert NICHT –
# settings.py liest os.getenv("SECRET_KEY"), nicht "DJANGO_SECRET_KEY".
# WARNUNG: Environment="DATABASE_NAME=..." funktioniert NICHT –
# settings.py liest DB_NAME, nicht DATABASE_NAME.

# Gunicorn mit 4 Workern
ExecStartPre=/bin/mkdir -p /path/to/your/project/logs
ExecStart=/path/to/your/project/venv/bin/gunicorn \
    --workers 4 \
    --bind unix:/path/to/your/project/homegym.sock \
    --timeout 180 \
    --access-logfile /path/to/your/project/logs/gunicorn-access.log \
    --error-logfile /path/to/your/project/logs/gunicorn-error.log \
    --log-level info \
    config.wsgi:application
# Hinweis zu --timeout 180:
# LLM-Calls an OpenRouter (70B-Modelle) dauern 60-120s.
# Gunicorn-Timeout MUSS höher sein als der längste API-Call.
# Nginx proxy_read_timeout muss ebenfalls auf 180s gesetzt sein.

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
